# gitlab-backup.service
[Unit]
Description=GitLab backup (application data + config)
After=gitlab.service

[Service]
Type=oneshot
# Run GitLab's built-in backup tool inside the running container.
# STRATEGY=copy prevents corruption if files change during archiving.
# The backup tar is written to /var/opt/gitlab/backups/ inside the container
# (= %h/data/backups/ on the host).
ExecStart=/usr/bin/podman exec systemd-gitlab gitlab-backup create STRATEGY=copy
# Sync the backup archive(s) to the staging area.
# podman unshare is required: the backup tar is created by the container's git user,
# which maps to a subuid on the host that the gitlab host user cannot access directly.
# Running inside the user namespace makes those UIDs visible as effective root.
ExecStart=/usr/bin/podman unshare /usr/bin/rsync -a --no-o --no-g --delete --chmod=Dg+rx,Fg+r %h/data/backups/ /var/backups/gitlab/backups/
# Backup config and secrets separately — gitlab-backup does not include /etc/gitlab/
ExecStart=/usr/bin/rsync -a --delete --chmod=Dg+rx,Fg+r %h/config/ /var/backups/gitlab/config/
